# slashjp/tags/2.5.0.233-4/debian/slash-backend.cron.d
*/5 *	* * *	root	/etc/init.d/slash keepalive
